# deploy/applyhome-alert.crontab
0 9 * * * /Users/osori/workbench/applyhome-alert-bot/scripts/run_alert.sh # applyhome-alert-bot
